# claude-telegram.service
[Unit]
Description=Claude Telegram Bot
After=network.target

[Service]
Type=simple
User=ariel
WorkingDirectory=/home/ariel/Development/claude-telegram
EnvironmentFile=/home/ariel/Development/claude-telegram/.env
ExecStart=/usr/bin/node --require /home/ariel/Development/claude-telegram/node_modules/tsx/dist/preflight.cjs --import file:///home/ariel/Development/claude-telegram/node_modules/tsx/dist/loader.mjs /home/ariel/Development/claude-telegram/index.ts
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-telegram

[Install]
WantedBy=multi-user.target
